# linux-backup/linux-home-backup.timer
[Unit]
Description=Run linux-home-backup daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target